# deploy/systemd/ups-pi-node.service
[Unit]
Description=ups-pi-node web service
After=network.target ups-pi-node-helper.service
Wants=ups-pi-node-helper.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/usr/lib/ups-pi-node
Environment="PATH=/var/lib/ups-pi-node/.venv/bin"
EnvironmentFile=-/etc/default/ups-pi-node
ExecStart=/var/lib/ups-pi-node/.venv/bin/python -m gunicorn --config gunicorn.conf.py wsgi:app
Restart=always

[Install]
WantedBy=multi-user.target
