# debian/prerm
#!/bin/bash
set -e

case "$1" in
    remove|upgrade)
        systemctl stop flask_site.service || true
        systemctl stop rpi2w-system-helper.service || true
        systemctl stop rpi2w-hotspot-fallback.timer || true
        systemctl disable flask_site.service || true
        systemctl disable rpi2w-system-helper.service || true
        systemctl disable rpi2w-hotspot-fallback.timer || true
        rm -f /etc/nginx/sites-enabled/rpi2w-ups-portal
        systemctl reload nginx || true
        ;;
esac

exit 0
